# grub-btrfsd.initd
#!/sbin/openrc-run
# Copyright 2021 Pascal Jaeger
# Distributed under the terms of the GNU General Public License v3

name="grub-btrfs daemon"
command="/usr/bin/grub-btrfs-openrc"
command_args="${snapshots} ${timeshift_auto}"
pidfile="/run/{RC_SVCNAME}.pid"
command_background=true

depend() {
   use localmount
}
