# grub-btrfs.path
[Unit]
Description=Monitors for new snapshots
DefaultDependencies=no
Requires=\x2esnapshots.mount
After=\x2esnapshots.mount
BindsTo=\x2esnapshots.mount

[Path]
PathModified=/.snapshots

[Install]
WantedBy=\x2esnapshots.mount
